# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/app/go.sh
