# supervisord.conf
[supervisord]
nodaemon=true

[program:mvn]
command=/app/go.sh
stopasgroup=true
